# oam/install_scripts/columnstore.service.in
[Unit]
Description=MariaDB Columnstore

[Service]
Type=oneshot
ExecStart=@ENGINE_BINDIR@/columnstore-post-install

# This service shall be considered active after start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
