# oam/install_scripts/controllernode.service.in
[Unit]
Description=controllernode
PartOf=columnstore.service
After=columnstore.service

[Service]
ExecStart=@ENGINE_BINDIR@/controllernode &> /var/log/mariadb/columnstore/controllernode.log
Restart=on-failure

[Install]
WantedBy=columnstore.service
