# oam/install_scripts/ddlproc.service.in
[Unit]
Description=DDLProc
PartOf=columnstore.service
After=columnstore.service

[Service]
ExecStart=@ENGINE_BINDIR@/DDLProc &> /var/log/mariadb/columnstore/ddlproc.log
Restart=on-failure


[Install]
WantedBy=columnstore.service
