# oam/install_scripts/dmlproc.service.in
[Unit]
Description=DMLProc
PartOf=columnstore.service
After=columnstore.service

[Service]
ExecStart=@ENGINE_BINDIR@/DMLProc
Restart=on-failure

[Install]
WantedBy=columnstore.service
