# oam/install_scripts/mariadb-columnstore.service.in
[Unit]
Description=mariadb-columnstore
After=network.target
Requires=mariadb-columnstore-cmapi.service

[Service]
Type=oneshot

# Execute a dummy program
ExecStart=@ENGINE_BINDIR@/mcs cluster status
ExecStop=@ENGINE_BINDIR@/mcs cluster stop

# This service shall be considered active after start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
