# oam/install_scripts/mcs-writeengineserver.service.in
[Unit]
Description=WriteEngineServer

# restart/stop mcs-writeengineserver on restart/stop of mcs-exemgr
PartOf=mcs-exemgr.service
After=network.target mcs-exemgr.service

[Service]
Type=simple

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNPROC=65536

ExecStart=/usr/bin/env bash -c "ldconfig -p | grep -m1 libjemalloc > /dev/null && LD_PRELOAD=$(ldconfig -p | grep -m1 libjemalloc | awk '{print $1}') exec @ENGINE_BINDIR@/WriteEngineServer || [ -f '@ENGINE_DATADIR@/libjemalloc.so.2' ] && LD_PRELOAD=@ENGINE_DATADIR@/libjemalloc.so.2 exec @ENGINE_BINDIR@/WriteEngineServer || exec @ENGINE_BINDIR@/WriteEngineServer"

Restart=on-failure
TimeoutStopSec=20
