# oam/install_scripts/mcs-primproc.service.in
[Unit]
Description=mcs-primproc

# restart/stop mcs-primproc on restart/stop of mcs-workernode or mcs-controllernode
PartOf=mcs-workernode@1.service mcs-workernode@2.service
PartOf=mcs-controllernode.service
After=network.target mcs-workernode@1.service mcs-workernode@2.service mcs-controllernode.service

[Service]
Type=forking

User=@DEFAULT_USER@
Group=@DEFAULT_GROUP@
LimitNOFILE=65536
LimitNPROC=65536

ExecStartPre=/usr/bin/env bash -c "ldconfig -p | grep -m1 libjemalloc > /dev/null || echo 'Please install jemalloc to avoid ColumnStore performance degradation and unexpected service interruptions.'"
ExecStart=/usr/bin/env bash -c "LD_PRELOAD=$(ldconfig -p | grep -m1 libjemalloc | awk '{print $1}') exec @ENGINE_BINDIR@/PrimProc"

Restart=on-failure
TimeoutStopSec=2
