# oam/install_scripts/mcs-primproc.service.in
[Unit]
Description=mcs-primproc
PartOf=mcs-workernode.service
PartOf=mcs-controllernode.service
After=mcs-controllernode.service

[Service]
Type=simple
Environment="SKIP_OAM_INIT=1"
ExecStart=@ENGINE_BINDIR@/PrimProc
Restart=on-failure
TimeoutStopSec=2

[Install]
WantedBy=mariadb-columnstore.service
