# oam/install_scripts/mcs-storagemanager.service.in
[Unit]
Description=storagemanager

[Service]
Type=forking

User=@DEFAULT_USER@
Group=@DEFAULT_GROUP@
LimitNOFILE=65536
LimitNPROC=65536

ExecStart=/usr/bin/env bash -c "LD_PRELOAD=$(ldconfig -p | grep -m1 libjemalloc | awk '{print $1}') exec @ENGINE_BINDIR@/StorageManager"
SuccessExitStatus=255
