# oam/install_scripts/mcs-ddlproc.service.in
[Unit]
Description=mcs-ddlproc
PartOf=mcs-exemgr.service
After=mcs-dmlproc.service

[Service]
Type=simple
Environment="SKIP_OAM_INIT=1"
ExecStart=@ENGINE_BINDIR@/DDLProc
Restart=on-failure
TimeoutStopSec=2

[Install]
WantedBy=mariadb-columnstore.service
