# oam/install_scripts/mcs-workernode.service.in
[Unit]
Description=mcs-workernode
PartOf=mariadb-columnstore.service
After=mcs-loadbrm.service

[Service]
Type=forking
ExecStart=@ENGINE_BINDIR@/workernode DBRM_Worker1
Restart=on-failure
ExecStop=-@ENGINE_BINDIR@/save_brm
ExecStopPost=/usr/bin/env bash -c "clearShm > /dev/null 2>&1"
TimeoutStopSec=2

[Install]
WantedBy=mariadb-columnstore.service
