# oam/install_scripts/mcs-dmlproc.service.in
[Unit]
Description=mcs-dmlproc

#restart/start mcs-dmlproc on mcs-controllernode restart/start
PartOf=mcs-controllernode.service 
After=network.target mcs-controllernode.service

[Service]
Type=forking

User=@DEFAULT_USER@
Group=@DEFAULT_GROUP@
LimitNOFILE=65536
LimitNPROC=65536

ExecStart=@ENGINE_BINDIR@/DMLProc

Restart=on-failure
TimeoutStopSec=2
