# oam/install_scripts/mcs-loadbrm.service.in
[Unit]
Description=loadbrm
PartOf=mariadb-columnstore.service
After=mariadb-columnstore.service
ConditionPathExists=/var/lib/columnstore/data1/systemFiles/dbrm/BRM_saves_current

[Service]
Type=simple
Environment="SKIP_OAM_INIT=1"
ExecStart=/usr/bin/bash -c "/usr/bin/load_brm /var/lib/columnstore/data1/systemFiles/dbrm/$(cat /var/lib/columnstore/data1/systemFiles/dbrm/BRM_saves_current)"

[Install]
WantedBy=mariadb-columnstore.service
