# oam/install_scripts/mcs-ddlproc.service.in
[Unit]
Description=mcs-ddlproc
PartOf=mcs-writeengineserver.service
After=mcs-dmlproc.service

[Service]
Type=simple
ExecStart=@ENGINE_BINDIR@/DDLProc
Restart=on-failure
TimeoutStopSec=2
