# oam/install_scripts/mcs-controllernode.service.in
[Unit]
Description=mcs-controllernode
After=network.target

[Service]
Type=forking

User=@DEFAULT_USER@
Group=@DEFAULT_GROUP@
LimitNOFILE=65536
LimitNPROC=65536

ExecStart=@ENGINE_BINDIR@/controllernode
ExecStop=@ENGINE_BINDIR@/mcs-stop-controllernode.sh $MAINPID

Restart=on-failure
