# oam/install_scripts/mcs-dmlproc.service.in
[Unit]
Description=mcs-dmlproc
PartOf=mcs-exemgr.service
After=mcs-exemgr.service

[Service]
Type=simple
ExecStart=@ENGINE_BINDIR@/DMLProc
Restart=on-failure
TimeoutStopSec=2

[Install]
WantedBy=mariadb-columnstore.service
