# oam/install_scripts/mcs-loadbrm.service.in
[Unit]
Description=loadbrm
PartOf=mcs-workernode.service
Before=mcs-workernode.service
ConditionPathExists=/var/lib/columnstore/data1/systemFiles/dbrm/BRM_saves_current

[Service]
Type=simple
ExecStart=/usr/bin/env bash -c "/usr/bin/load_brm /var/lib/columnstore/data1/systemFiles/dbrm/$(cat /var/lib/columnstore/data1/systemFiles/dbrm/BRM_saves_current)"

[Install]
WantedBy=mariadb-columnstore.service
WantedBy=mcs-workernode.service
