# oam/install_scripts/mcs-controllernode.service.in
[Unit]
Description=mcs-controllernode
PartOf=mcs-workernode.service
After=mcs-workernode.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 7
ExecStart=@ENGINE_BINDIR@/controllernode fg
Restart=on-failure
ExecStop=@ENGINE_BINDIR@/mcs-stop-controllernode.sh $MAINPID
