# oam/install_scripts/mcs-exemgr.service.in
[Unit]
Description=mcs-exemgr
PartOf=mcs-writeengineserver.service
After=mcs-writeengineserver.service

[Service]
Type=simple
Environment="SKIP_OAM_INIT=1"
ExecStart=@ENGINE_BINDIR@/ExeMgr
Restart=on-failure
ExecStop=/usr/bin/env bash -c "/bin/kill -15 $(pidof ExeMgr) && sleep 2 && /bin/kill -9 $(pidof ExeMgr) || /bin/true"

[Install]
WantedBy=mariadb-columnstore.service
