# oam/install_scripts/mcs-controllernode.service.in
[Unit]
Description=mcs-controllernode
PartOf=mcs-workernode.service
After=mcs-workernode.service

[Service]
Type=forking
Environment="SKIP_OAM_INIT=1"
ExecStart=@ENGINE_BINDIR@/controllernode
Restart=on-failure
ExecStop=@ENGINE_BINDIR@/mcs-stop-controllernode.sh
TimeoutStopSec=15

[Install]
WantedBy=mariadb-columnstore.service
