# oam/install_scripts/mcs-primproc.service.in
[Unit]
Description=mcs-primproc

# restart/stop mcs-primproc on restart/stop of mcs-workernode or mcs-controllernode
PartOf=mcs-workernode.service
PartOf=mcs-controllernode.service
After=network.target mcs-workernode.service mcs-controllernode.service

[Service]
Type=simple

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNPROC=65536

#ExecStartPre=/usr/bin/env bash -c "ldconfig -p | grep -m1 libjemalloc > /dev/null || echo 'Please install jemalloc to avoid ColumnStore performance degradation and unexpected service interruptions.'"
ExecStart=/usr/bin/env bash -c "LD_PRELOAD=@ENGINE_DATADIR@/libjemalloc.so.2 exec @ENGINE_BINDIR@/PrimProc"
ExecStartPost=/bin/sleep 2

Restart=on-failure
TimeoutStopSec=2
