# oam/install_scripts/mcs-loadbrm.service.in
[Unit]
Description=loadbrm
PartOf=mcs-workernode.service
After=mcs-storagemanager.service

[Service]
Type=simple
Environment="SKIP_OAM_INIT=1"
ExecStart=@ENGINE_BINDIR@/mcs-loadbrm.py

[Install]
WantedBy=mariadb-columnstore.service
WantedBy=mcs-workernode.service
