# oam/install_scripts/mcs-dmlproc.service.in
[Unit]
Description=mcs-dmlproc
PartOf=mcs-writeengineserver.service
After=mcs-writeengineserver.service

[Service]
Type=simple
ExecStart=/usr/bin/env bash -c "/bin/sleep 2 && @ENGINE_BINDIR@/DMLProc"
Restart=on-failure
TimeoutStopSec=2
