# oam/install_scripts/mcs-controllernode.service.in
[Unit]
Description=mcs-controllernode
PartOf=mcs-workernode.service
After=mcs-workernode.service

[Service]
Type=forking
ExecStartPre=/usr/bin/env bash -c "systemctl start mcs-workernode"
ExecStart=@ENGINE_BINDIR@/controllernode
Restart=on-failure
ExecStop=@ENGINE_BINDIR@/mcs-stop-controllernode.sh $MAINPID
