# oam/install_scripts/mcs-workernode.service.in
[Unit]
Description=mcs-workernode
After=mcs-loadbrm.service

[Service]
Type=simple
ExecStart=@ENGINE_BINDIR@/workernode DBRM_Worker1 fg
Restart=on-failure
ExecStop=/usr/bin/env bash -c "kill -15 $MAINPID"
ExecStopPost=-@ENGINE_BINDIR@/save_brm
ExecStopPost=/usr/bin/env bash -c "clearShm > /dev/null 2>&1"
TimeoutStopSec=5
