# oam/install_scripts/mcs-writeengineserver.service.in
[Unit]
Description=WriteEngineServer
PartOf=mcs-primproc.service
After=mcs-primproc.service

[Service]
Type=simple
Environment="LD_PRELOAD=libjemalloc.so.1 libjemalloc.so.2"
ExecStart=@ENGINE_BINDIR@/WriteEngineServer
Restart=on-failure
TimeoutStopSec=2

[Install]
WantedBy=mariadb-columnstore.service
