# oam/install_scripts/mcs-controllernode.service.in
[Unit]
Description=mcs-controllernode
PartOf=mariadb-columnstore.service
After=mariadb-columnstore.service

[Service]
Type=forking
Environment="SKIP_OAM_INIT=1"
ExecStart=@ENGINE_BINDIR@/controllernode
Restart=on-failure
ExecStop=/bin/kill -9 $MAINPID

[Install]
WantedBy=mariadb-columnstore.service
