# oam/install_scripts/mcs-storagemanager.service.in
[Unit]
Description=storagemanager

[Service]
Type=simple

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNPROC=65536

ExecStart=/usr/bin/env bash -c "LD_PRELOAD=@ENGINE_DATADIR@/libjemalloc.so.2 exec @ENGINE_BINDIR@/StorageManager"
SuccessExitStatus=255
