# oam/install_scripts/mcs-writeengineserver.service.in
[Unit]
Description=WriteEngineServer

# restart/stop mcs-writeengineserver on restart/stop of mcs-primproc
# PartOf=mcs-primproc.service
After=network.target mcs-primproc.service

[Service]
Type=forking

User=@DEFAULT_USER@
Group=@DEFAULT_GROUP@
LimitNOFILE=65536
LimitNPROC=65536
LimitCORE=@CORE_DUMPS@

ExecStart=/usr/bin/env bash -c "@WRITEENGINE_ALLOC_CONFIG@ @LD_PRELOAD_STRING@ exec @ENGINE_BINDIR@/WriteEngineServer"

Restart=on-failure
TimeoutStopSec=20
