# oam/install_scripts/mcs-loadbrm.service.in
[Unit]
Description=loadbrm
PartOf=mcs-workernode.service
After=mcs-storagemanager.service

[Service]
Type=simple
ExecStart=@ENGINE_BINDIR@/mcs-loadbrm.py

[Install]
WantedBy=mcs-workernode.service
