# oam/install_scripts/mcs-storagemanager.service.in
[Unit]
Description=storagemanager

[Service]
Type=forking

User=@DEFAULT_USER@
Group=@DEFAULT_GROUP@
LimitNOFILE=65536
LimitNPROC=65536
LimitCORE=@CORE_DUMPS@

ExecStart=/usr/bin/env bash -c "@STORAGEMANAGER_ALLOC_CONFIG@ @LD_PRELOAD_STRING@ exec @ENGINE_BINDIR@/StorageManager"
SuccessExitStatus=255
