# oam/install_scripts/primproc.service.in
[Unit]
Description=PrimProc
PartOf=columnstore.service
After=columnstore.service

[Service]
ExecStart=@ENGINE_BINDIR@/PrimProc &> /var/log/mariadb/columnstore/primproc.log
Restart=on-failure

[Install]
WantedBy=columnstore.service
