# oam/install_scripts/workernode.service.in
[Unit]
Description=workernode
PartOf=columnstore.service
After=columnstore.service

[Service]
ExecStartPre=@ENGINE_BINDIR@/bash -c "clearShm > /dev/null 2>&1"
ExecStart=@ENGINE_BINDIR@/workernode DBRM_Worker1 &> /var/log/mariadb/columnstore/workernode.log
Restart=on-failure

[Install]
WantedBy=columnstore.service
