# oam/install_scripts/writeengine.service.in
[Unit]
Description=WriteEngineServer
PartOf=columnstore.service
After=columnstore.service

[Service]
ExecStart=/usr/bin/WriteEngineServer
Restart=on-failure

[Install]
WantedBy=columnstore.service
